# deploy/vehicle-counter.service
[Unit]
Description=jetson-yolo-cpp vehicle counter
Requires=docker.service
After=docker.service

[Service]
Type=oneshot
RemainAfterExit=yes
WorkingDirectory=/opt/jetson-yolo-cpp
ExecStart=/usr/bin/docker compose -f docker/docker-compose.yml up -d
ExecStop=/usr/bin/docker compose -f docker/docker-compose.yml down
TimeoutStartSec=120

[Install]
WantedBy=multi-user.target
